# recipes-multimedia/pulseaudio/files/pulseaudio.service
[Unit]
Description=PulseAudio Sound System
After=alsa-restore.service

[Service]
BusName=org.pulseaudio.Server
ExecStart=/usr/bin/pulseaudio --system
Restart=always

[Install]
WantedBy=multi-user.target
